# configs/oai-logrotate.service
[Unit]
Description=OAI hourly logrotate
Documentation=man:logrotate(8)

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate /etc/logrotate.conf
